# agents/MoxWiki/setup/moxwiki.conf
description "Upstart file for running the MoxWiki agent on system boot"

start on runlevel [2345]
stop on shutdown

exec /srv/mox/agents/MoxWiki/moxwiki/app.py

