# python_agents/setup/mox-elk-log.in.conf
description "Upstart file for running the log integration agent on system boot"

start on runlevel [2345]
stop on shutdown

exec %PATH%/mox_elk_log.py

