# config/etc/init/moxagent.conf
description "Upstart file for running the oio_api on system boot"

start on runlevel [2345]
stop on shutdown

exec /srv/mox/agent/agent.sh

